# default_config/network/if-up.d/berryio_email_ip
#!/bin/bash
# Script to email a link to the BerryIO web interface
berryio email_ip
